# pi/hatel.service
[Unit]
Description=Hatel in systemd.

[Service]
ExecStart=/home/alarm/hatel pi

[Install]
WantedBy=multi-user.target
